# backend/api_dashboard.service
[Unit]
Description=API Usage Dashboard Backend
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/share/api-dashboard/venv/bin/python %h/.local/share/api-dashboard/api_dashboard_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
